# config/systemd/user/checkupdates.timer
[Unit]
Description=check arch repository for updates

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
