# config/systemd/user/fetchrss.timer
[Unit]
Description=fetch new articles from feeds

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h 

[Install]
WantedBy=timers.target
